# organizations/palermo/etc/uci-defaults/902_ninux-olsrd
#!/bin/sh
uci batch <<-EOF
    # mesh
    set olsrd.ninux=Interface
    set olsrd.ninux.ignore=0
    set olsrd.ninux.interface='ninux'
    set olsrd.ninux.Mode=mesh
    commit olsrd
EOF
exit 0
